# one-api.service
# File path: /etc/systemd/system/one-api.service
# sudo systemctl daemon-reload
# sudo systemctl start one-api
# sudo systemctl enable one-api
# sudo systemctl status one-api
[Unit]
Description=One API Service
After=network.target

[Service]
User=ubuntu  # Make sure to change the username
WorkingDirectory=/path/to/one-api  # Make sure to change the path
ExecStart=/path/to/one-api/one-api --port 3000 --log-dir /path/to/one-api/logs  # Make sure to change the path and port
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
